# cmder/config/profile.d/ls.alias.ps1
# ls
# Pre: Install Get-ChildItemColor
# Install-Module Get-ChildItemColor (or Install-Module Get-ChildItemColor -Scope CurrentUser
# Get-ChildItemColor Homepage: https://github.com/joonro/Get-ChildItemColor

Set-Alias ls Get-ChildItemColorFormatWide -option AllScope